# launchd/com.user.claude-sessions-sync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.user.claude-sessions-sync — WatchPaths on the primary Claude Desktop
  userData's `Local State` file (Chromium rewrites it on app-close). On
  trigger, runs claude-sessions-sync.sh which one-way rsyncs IndexedDB /
  Local Storage / Session Storage from the primary to every configured mirror.

  Template placeholders (substituted at install time by lib/install-launchd.sh):
    __SYNC_BIN__            — full path to bin/claude-sessions-sync.sh
    __PRIMARY_LOCALSTATE__  — full path to <primary userData>/Local State
    __LOG_PATH__            — full path to the WorkingDirectory-relative log

  ProcessType=Background: nice-value 10, low-I/O priority (never contends with
  active user work). LowPriorityIO=true also disables burst I/O prioritisation.

  Nice=10 + ThrottleInterval=120 prevents thrash if Chromium touches Local State
  in bursts. 5-second grace inside the script gives Chromium time to flush.

  StartInterval=21600 (6h) is a safety net for the case where the primary crashes
  and Local State never gets rewritten — WatchPaths would miss the event forever.
  The sync worker hard-skips when any Claude Desktop is running, so a periodic
  wake while both apps are open is a cheap no-op.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-sessions-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SYNC_BIN__</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>__PRIMARY_LOCALSTATE__</string>
    </array>

    <key>StandardOutPath</key>
    <string>__LOG_PATH__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_PATH__</string>

    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>
    <key>LowPriorityIO</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>120</integer>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
