# launchd/com.user.claude-metadata-symlink.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.user.claude-metadata-symlink — WatchPaths on every configured mirror's
  <userData>/claude-code-sessions, <userData>/local-agent-mode-sessions, and
  <userData>/config.json. The dirs get an mtime bump the moment Claude Desktop
  creates the per-account UUID subdir on Code first-use. config.json gets an
  mtime bump on OAuth sign-in itself (Desktop caches tokens + writes
  lastKnownAccountUuid). Watching all three lets us install the metadata
  symlinks the moment the user signs in, without any "open Code once inside
  Desktop then quit" manual step.

  On trigger, runs bin/claude-metadata-watcher.sh which iterates every mirror
  and installs any missing per-account-folder symlinks so the mirror's desktop
  UI resolves to the primary's session store. UUIDs are resolved via on-disk
  find when the dirs exist, and via config.json (lastKnownAccountUuid +
  newest-timestamp dxt:allowlistLastUpdated:<orgUuid>) otherwise — see the
  header comment in lib/metadata-symlinks.sh for detail.

  Template placeholders (substituted at install time by lib/install-launchd.sh):
    __WATCH_BIN__            — full path to bin/claude-metadata-watcher.sh under
                               the libexec copy (~/.local/libexec/claude-
                               multiacct/bin), NOT the repo checkout — a
                               background launchd agent cannot exec from a
                               TCC-protected container like ~/Documents.
    __WATCHPATHS_ENTRIES__   — one <string>…</string> per watched path (three
                               per mirror: claude-code-sessions, local-agent-
                               mode-sessions, and config.json). Injected via
                               `sed r file; d`.
    __CMA_AUTO_RESTART__     — "0" (default) or "1". Gates the watcher's auto
                               quit+relaunch of a mirror on first symlink
                               install. Symlink install + sentinel always
                               happen; only the disruptive kill/relaunch of a
                               live mirror is opt-in. Default OFF.
    __LOG_PATH__             — full path to the log file

  ProcessType=Background + Nice=10 + LowPriorityIO=true: install-time symlink
  work is tiny (<50ms) but keep it out of foreground user-work I/O anyway.

  ThrottleInterval=15 coalesces the burst of dir mutations Desktop performs
  while setting up the account+org UUID subtree on sign-in.

  RunAtLoad=false: reloading the agent doesn't imply mirrors changed state;
  we only care about mtime-change events under the WatchPaths.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-metadata-symlink</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WATCH_BIN__</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CMA_AUTO_RESTART</key>
        <string>__CMA_AUTO_RESTART__</string>
    </dict>

    <key>WatchPaths</key>
    <array>
__WATCHPATHS_ENTRIES__
    </array>

    <key>StandardOutPath</key>
    <string>__LOG_PATH__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_PATH__</string>

    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>
    <key>LowPriorityIO</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
