# launchd/com.user.claude-primary-patch-refresh.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.user.claude-primary-patch-refresh — WatchPaths on the PRIMARY Claude
  Desktop bundle's Info.plist. Squirrel rewrites this file at the end of every
  Claude Desktop update, which reverts /Applications/Claude.app to Anthropic's
  pristine (unpatched, Anthropic-signed) shape. On mtime change, this agent
  re-runs bin/claude-primary-patch-refresh.sh which invokes `claude-multiacct
  primary-patch` — re-injecting the propagation + rc-enforcer IIFEs, rewriting
  the managed-config plist route, and re-signing inside-out with the loose
  Designated Requirement so the primary is back in the patched-but-updatable
  state the user asked for.

  This agent is load-bearing precisely BECAUSE the loose-DR re-sign keeps the
  primary Squirrel-updatable: updates now actually apply, so Squirrel drops a
  fresh Anthropic-signed bundle regularly, and this agent re-patches it.

  Load-bearing ordering vs. the clone-refresh agent:
    Squirrel touches Info.plist → BOTH agents fire in parallel:
      1. primary-patch-refresh  → patches the primary in place (~30 s)
      2. clone-refresh          → ditto's the primary → mirror clones (~30 s)
    The two share a lock file at /tmp/claude-multiacct-refresh.lock so
    primary-patch-refresh always completes BEFORE clone-refresh's ditto
    captures the primary bundle. That ordering means clone-refresh always
    ditto's the ALREADY-PATCHED primary — the mirror clone's own asar-patch
    step then no-ops on the "already patched" idempotency check.

  Both agents include the same Info.plist WatchPath: the trigger event is the
  same, but each script guards its own critical section with the shared flock
  so they serialize even under a concurrent-launchd fire.

  Template placeholders (substituted at install time by lib/install-launchd.sh):
    __PATCH_BIN__      — full path to bin/claude-primary-patch-refresh.sh
                         (rendered to the libexec copy, NOT the repo checkout)
    __WATCH_PATH__     — full path to <primary Claude.app>/Contents/Info.plist
    __LOG_PATH__       — full path to the log file

  ProcessType=Background + Nice=10 + LowPriorityIO=true: the patch step is
  asar extract + repack + inside-out codesign — I/O-heavy. Keep out of the way
  of interactive user work.

  ThrottleInterval=60 prevents a burst if Squirrel touches the plist multiple
  times during an update.

  RunAtLoad=false: reloading the agent doesn't imply Squirrel ran; only fire
  on real Info.plist mutation.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-primary-patch-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PATCH_BIN__</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>__WATCH_PATH__</string>
    </array>

    <key>StandardOutPath</key>
    <string>__LOG_PATH__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_PATH__</string>

    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>
    <key>LowPriorityIO</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
