# launchd/com.user.claude-clone-refresh.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.user.claude-clone-refresh — WatchPaths on the PRIMARY Claude Desktop
  bundle's Info.plist. Squirrel (the Electron autoupdater) rewrites this
  file at the end of every Claude Desktop update; the mtime change fires
  this agent. On trigger, runs bin/claude-clone-refresh.sh which rebuilds
  every mirror instance's clone from the new /Applications/Claude.app so
  the mirrors track the primary's version.

  Without this agent, mirror clones drift behind the primary after every
  auto-update — they'd keep launching the previous Claude Desktop version
  until the user manually ran `claude-multiacct refresh-clones`.

  Template placeholders (substituted at install time by lib/install-launchd.sh):
    __REFRESH_BIN__    — full path to bin/claude-clone-refresh.sh
    __WATCH_PATH__     — full path to <primary Claude.app>/Contents/Info.plist
    __LOG_PATH__       — full path to the log file

  ProcessType=Background + Nice=10 + LowPriorityIO=true: rebuilds involve
  a 745 MB `ditto` + a full codesign over Electron helpers, both I/O-heavy.
  Keep them out of the way of interactive user work.

  ThrottleInterval=60 prevents a burst of rebuilds if Squirrel touches the
  plist multiple times during an update.

  RunAtLoad=false: we don't want to rebuild every clone on every login just
  because the agent was reloaded; only rebuild when Info.plist actually changes.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-clone-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REFRESH_BIN__</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>__WATCH_PATH__</string>
    </array>

    <key>StandardOutPath</key>
    <string>__LOG_PATH__</string>
    <key>StandardErrorPath</key>
    <string>__LOG_PATH__</string>

    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>10</integer>
    <key>LowPriorityIO</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
